# systemd/voice-typing.service
[Unit]
Description=Local voice typing daemon (RealtimeSTT)
After=pipewire.service ydotool.service

[Service]
# Import the Wayland/X display variables from the systemd USER manager environment into THIS
# service's process environment. systemd user services do NOT automatically inherit arbitrary
# variables from `systemctl --user show-environment`; only a curated set is passed through.
# `wtype` (the default typing backend, PRD §4.3) is a Wayland `virtual-keyboard-v1` CLIENT and
# cannot connect to the compositor without WAYLAND_DISPLAY — so without this import every
# finalized utterance silently fails wtype and falls back to ydotool (~1.2s extra latency + the
# PRD-flagged non-ASCII/layout quirks). DISPLAY covers XWayland-focused windows.
#
# The user manager's environment is populated at graphical-session launch (Hyprland exports
# WAYLAND_DISPLAY/DISPLAY via the session; `systemctl --user import-environment` is the standard
# way a compositor seeds it). import-environment here is idempotent and a no-op if the vars are
# already set in the unit's environment. Runs as ExecStartPre so a failure (var not set yet) is
# surfaced in the journal before the daemon starts rather than as per-utterance wtype failures.
ExecStartPre=/usr/bin/systemctl --user import-environment WAYLAND_DISPLAY DISPLAY
# ExecStart points at the LD_LIBRARY_PATH wrapper (voice_typing/launch_daemon.sh), NOT python
# directly. faster-whisper / ctranslate2 load cuBLAS + cuDNN 9 shared objects from the nvidia-*-cu12
# wheels; cuDNN 9's split sub-libs (libcudnn_ops.so.9, libcudnn_cnn.so.9, ...) are resolved by the
# dynamic linker and the wheels ship no $ORIGIN RUNPATH for them, so they MUST be on LD_LIBRARY_PATH
# or the daemon dies at startup with "cannot load libcudnn_ops.so.9" (PRD §8 risk row #1).
#
# LD_LIBRARY_PATH is read by ld.so ONLY at process exec (execve) — mutating os.environ inside the
# already-running python has NO effect. The wrapper therefore exports it BEFORE exec'ing python.
#
# DO NOT add `Environment=LD_LIBRARY_PATH=...` here. The wrapper recomputes the lib dirs from the
# LIVE installed nvidia wheels on every launch, so this unit survives `uv sync` reinstalls and
# version bumps without edits. A baked value would go stale and break the service on the next sync.
# Debugging a cuDNN load failure: see the header comment in voice_typing/launch_daemon.sh
# (journalctl --user -u voice-typing -e; systemctl --user show voice-typing -p Environment;
#  LD_DEBUG=libs voice_typing/launch_daemon.sh 2>&1 | grep -i cudnn).
#
# The daemon boots NOT-LISTENING — it never hot-mics on boot / restart. `voicectl start` or
# `voicectl toggle` arms the mic over the control socket ($XDG_RUNTIME_DIR/voice-typing/control.sock).
# (PRD §4.9.) Verify with: voicectl status   ->   listening: off
ExecStart=/home/dustin/projects/voice-typing/voice_typing/launch_daemon.sh
Restart=on-failure
RestartSec=2
# Bound systemd's stop so the daemon is never SIGKILLed at the 90s DEFAULT (PRD §8 risk row:
# "recorder.shutdown() hangs ~90s … SIGKILL after systemd TimeoutStopSec"; journal:
# 'Failed with result timeout'). With no TimeoutStopSec, systemd applies its 90s default.
# The daemon's own _bounded_shutdown(timeout=5) (P1.M1.T1.S2 + P1.M1.T2.S2) makes
# VoiceTypingDaemon.shutdown() return within ~7s even when RealtimeSTT's shutdown_recorder()
# wedges at an unbounded threading.Thread.join() (root-caused in P1.M1.T1.S1). With daemon
# single-flight (P1.M1.T2.S1: exactly ONE teardown on the SIGTERM path) the 15s budget =
# 5s bounded teardown + 5s signal/shutdown coordination + 5s headroom.
# Normal stop exits in ~seconds; SIGKILL@15s is the outer safety net for a regression wedge
# (it skips the daemon's clean VRAM/process release, so bounded teardown remains the primary path).
TimeoutStopSec=15

[Install]
WantedBy=default.target
